# service/solana-sync/run
#!/bin/sh

app_base=/var/run/solana exec /sbin/solana-sync